# env/supervisor/worker/websockets.conf
[program:websockets]
command=/opt/bitnami/php/bin/php /var/www/websockets/artisan websockets:serve --port=6060
numprocs=1
autostart=true
autorestart=true
user=root